# llm/llama-server.service
[Unit]
Description=llama.cpp server for embedded AI assistant
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/llama-server -m /opt/ai-assistant/models/Qwen2.5-0.5B-Instruct-Q4_K_M.gguf --host 0.0.0.0 --port 11434 -c 2048 --threads 2 -ngl 0 --chat-template chatml
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
